# Mare/ExampleConfig/Systemd/mare.target
[Unit]
Description=Mare Synchronos
After=network.target

[Install]
WantedBy=multi-user.target
